# root/etc/profile.d/90_global_alias.sh.tmpl
#!/bin/sh

{{- if lookPath "bsdtar" -}}
alias tar=bsdtar
{{- end -}}
