# root/etc/profile.d/90_global_environment.sh.tmpl
#!/bin/sh

# sh prompt username:hostname:path$

export PS1='\u@\h:\w\$ '

# Set PATH, MANPATH, etc., for Homebrew.
eval "$(/home/linuxbrew/.linuxbrew/bin/brew shellenv)"
