# root/etc/profile.d/90_global_environment.sh.tmpl
#!/bin/sh

# sh prompt username:hostname:path$

export PS1='\u@\h:\w\$ '

{{ if lookPath "bsdtar" -}}
alias tar=bsdtar
{{- end }}
